# packages/.fluent-bit/linux/fluent-bit_3.2.3_linux_amd64/install.sh
#!/bin/sh

wget https://raw.githubusercontent.com/fluent/fluent-bit/master/install.sh | sh

systemctl disable fluent-bit
